# .claude/voice-server/macos-service/com.paivoice.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paivoice.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bun</string>
        <string>run</string>
        <string>/Users/daniel/.claude/voice-server/server.ts</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8888</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/daniel</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/daniel/.claude/voice-server</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/daniel/.claude/voice-server/logs/voice-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/daniel/.claude/voice-server/logs/voice-server-error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>